# lcs.service
[Unit]
Description=LCS Train Monitor Service
After=network.target mysql.service

[Service]
# User/Group set to root as requested (path is /root/LCS)
User=root
Group=root

# Path to where the repo is cloned
WorkingDirectory=/root/LCS

# Command to run the app
ExecStart=/usr/bin/python3 /root/LCS/lcs_monitor.py

# Restart automatically if it crashes
Restart=always
RestartSec=5

# Ensure output is logged immediately
Environment=PYTHONUNBUFFERED=1
# Environment="LCS_DB_HOST=sql.internal" # Example: Override config here if needed

[Install]
WantedBy=multi-user.target
